# data/connectivity-manager.service.in
[Unit]
Description=Connectivity Manager

[Service]
Type=dbus
BusName=com.luxoft.ConnectivityManager
ExecStart=@bindir@/connectivity-manager
