# data/template-dbus-service.service.in
[Unit]
Description=Template D-Bus service

[Service]
Type=dbus
BusName=com.luxoft.Template
ExecStart=@bindir@/template-dbus-service
